# scraper/condodork-scraper.timer
[Unit]
Description=Run CondoDork scraper every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
